# docker/php/services.d/php/run
#!/usr/bin/execlineb -P
/usr/bin/with-contenv
/usr/local/sbin/php-fpm -R
